# systemd/avatar-gateway.service
[Unit]
Description=RK3588 USB Avatar Gateway
After=network.target local-fs.target uvc-gadget-setup.service
Wants=uvc-gadget-setup.service

[Service]
Type=simple
WorkingDirectory=/opt/rk3588-avatar-gateway
EnvironmentFile=-/etc/default/avatar-gateway
Environment=OUTPUT_MODE=network
Environment=OUTPUT_DEVICE=/dev/video43
Environment=NETWORK_HOST=0.0.0.0
Environment=NETWORK_PORT=8080
Environment=NETWORK_PATH=/mjpeg
Environment=NETWORK_JPEG_QUALITY=85
Environment=FALLBACK_STYLE=normal
ExecStart=/opt/rk3588-avatar-gateway/scripts/run_avatar_gateway.sh
Restart=always
RestartSec=1
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target